# supervisor_conf/pan_tilt.conf
[program:pan_tilt]
command=/usr/local/bin/pan_tilt.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/pan_tilt.err.log
stdout_logfile=/var/log/pan_tilt.out.log
